# crontab
0 5 * * * ./blocksite enable
0 18 * * * ./blocksite disable
